# thermal.service
[Unit]
Description=Run Thremal Camera Programs.

[Service]
WorkingDirectory=/home/pi
ExecStart=/bin/bash /home/pi/systemdthermal.sh start
Restart=always
RestartSec=7
User=root

[Install]
WantedBy=multi-user.target
